# fs_mini_mdev/etc/inittab
console::askfirst:-/bin/sh
::sysinit:/etc/init.d/rcS
